# scheduler/systemd/user/cursor-edamame-extrapolator.service
[Unit]
Description=Cursor EDAMAME extrapolator
After=default.target

[Service]
Type=oneshot
WorkingDirectory=__PACKAGE_ROOT__
ExecStart=__NODE_BIN__ __PACKAGE_ROOT__/service/cursor_extrapolator.mjs --config __CONFIG_PATH__ --json
StandardOutput=append:__STATE_DIR__/cursor-extrapolator.log
StandardError=append:__STATE_DIR__/cursor-extrapolator.err
