# scheduler/systemd/user/cursor-edamame-extrapolator.timer
[Unit]
Description=Run Cursor EDAMAME extrapolator every 2 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=120s
Unit=cursor-edamame-extrapolator.service

[Install]
WantedBy=timers.target
