# scheduler/systemd/user/cursor-edamame-verdict-reader.timer
[Unit]
Description=Run Cursor EDAMAME verdict reader every 3 minutes

[Timer]
OnBootSec=60s
OnUnitActiveSec=180s
Unit=cursor-edamame-verdict-reader.service

[Install]
WantedBy=timers.target
